# Splunk/rules-emerging-threats/Execution/diamond_sleet_apt_scheduled_task_creation.spl
# Title: Diamond Sleet APT Scheduled Task Creation
# Author: Nasreddine Bencherchali (Nextron Systems)
# Date: 2023-10-24
# Level: critical
# Description: Detects registry event related to the creation of a scheduled task used by Diamond Sleet APT during exploitation of Team City CVE-2023-42793 vulnerability
# MITRE Tactic: Execution
# Tags: attack.execution, attack.privilege-escalation, attack.persistence, attack.t1053.005, detection.emerging-threats

EventID=4698 TaskContent="*uTYNkfKxHiZrx3KJ*" TaskName="\\Windows TeamCity Settings User Interface"