# Splunk/rules-emerging-threats/Privilege Escalation/suspicious_vbscript_un2452_pattern.spl
# Title: Suspicious VBScript UN2452 Pattern
# Author: Florian Roth (Nextron Systems)
# Date: 2021-03-05
# Level: high
# Description: Detects suspicious inline VBScript keywords as used by UNC2452
# MITRE Tactic: Privilege Escalation
# Tags: attack.privilege-escalation, attack.persistence, attack.t1547.001, detection.emerging-threats

CommandLine="*Execute*" CommandLine="*CreateObject*" CommandLine="*RegRead*" CommandLine="*window.close*" CommandLine="*\\Microsoft\\Windows\\CurrentVersion*" NOT CommandLine="*\\Software\\Microsoft\\Windows\\CurrentVersion\\Run*"