# Splunk/rules/linux/builtin/cron/modifying_crontab.spl
# Title: Modifying Crontab
# Author: Pawel Mazur
# Date: 2022-04-16
# Level: medium
# Description: Detects suspicious modification of crontab file.
# MITRE Tactic: Privilege Escalation
# Tags: attack.privilege-escalation, attack.execution, attack.persistence, attack.t1053.003
# False Positives:
#   - Legitimate modification of crontab

"REPLACE"